# release/src-ra-5010/linux/vendors/Mediatek/MT7621/inittab
::sysinit:/etc_ro/rcS
ttyS1::respawn:/bin/sh
